# build/templates/tot-backend.service.tmpl
[Unit]
Description=tot-backend
After=syslog.target
After=network.target
Requires=nginx.service
Requires=tot-system.service

[Service]
Type=simple
PIDFile=/run/tot-backend.pid
Environment=NODE_CONFIG_DIR=/etc/tot/backend/

User=tot
Group=tot

ExecStart=/usr/bin/node /usr/src/tot/backend/src/index.js
ExecStop=/bin/kill -HUP $MAINPID
ExecReload=/bin/kill -HUP $MAINPID
KillMode=process
Restart=always
TimeoutSec=300

[Install]
WantedBy=multi-user.target 